# rc.d/reboot
#!/bin/sh

###
# Vacteria System init script, Created on lun ene 10 14:12:34 CST 2011
# This file is free software and is released under GPL v2.1 or grant.
###

### BEGIN INIT INFO
# Provides:          reboot
# Required-Start:    
# Required-Stop:     
# Should-Start:      
# Should-Stop:       
# X-Start-Before:    
# X-Stop-After:      
# Default-Start:     6
# Default-Stop:      
# X-Interactive:     
# X-vct-Provided-By: Vacteria
# Short-Description: Reboots the system.
# Description:       Reboots the system.
### END INIT INFO

. /lib/init/init-functions

rc_start()
{
      rc_message "Restarting system"
      reboot -d -f -i
}

case ${1} in
	start   ) rc_start   ;;
	*       ) rc_error   ;;
esac
